# autostart/lxqt-notifications.desktop.in
[Desktop Entry]
Type=Application
Name=Notification Daemon
TryExec=lxqt-notificationd
Exec=lxqt-notificationd
OnlyShowIn=Razor;
X-LxQt-Module=true

#TRANSLATIONS_DIR=translations
