# autostart/lxqt-notifications.desktop.in
[Desktop Entry]
Type=Application
TryExec=lxqt-notificationd
Exec=lxqt-notificationd
OnlyShowIn=LXQt;
X-LXQt-Module=true

#TRANSLATIONS_DIR=translations
